# ops/systemd/liminal-living-site.service.example
[Unit]
Description=Liminal Living Site Evolution Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=liminal
Group=liminal
WorkingDirectory=/opt/liminal
Environment=NODE_ENV=production
EnvironmentFile=/etc/liminal/living-site.env
ExecStart=/usr/bin/env node /opt/liminal/bin/sinter site evolve
Restart=on-failure
RestartSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/www/kyanitelabs/liminal-asset /var/lib/liminal /tmp

[Install]
WantedBy=multi-user.target
